# deploy/corpusagent2-stack.service
[Unit]
Description=CorpusAgent2 Docker stack
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/dongtten/corpusagent2/deploy
EnvironmentFile=/home/dongtten/corpusagent2/.env
ExecStart=/bin/sh -lc '/usr/bin/docker compose -f docker-compose.yml -f docker-compose.mcp.yml up -d --no-recreate postgres opensearch && /usr/bin/docker compose -f docker-compose.yml -f docker-compose.mcp.yml up -d --build --no-deps corpusagent2-api corpusagent2-mcp'
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.mcp.yml stop corpusagent2-api corpusagent2-mcp
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
